# deploy/hestia-agent.service
[Unit]
Description=Hestia agent (Telegram + Scheduler)
After=network.target hestia-llama.service
Requires=hestia-llama.service

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/hestia
ExecStart=/opt/hestia/.venv/bin/hestia --config /opt/hestia/config.py telegram
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
